# garage-management.service
[Unit]
Description=Garage Management System
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/garage-management-system
Environment=PATH=/usr/bin:/usr/local/bin
Environment=FLASK_ENV=production
Environment=PORT=5001
ExecStart=/usr/bin/python3 src/main.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=garage-management

[Install]
WantedBy=multi-user.target
